# docker/supervisor-websockets.conf
[program:websockets]
process_name=%(program_name)s
command=php /var/www/html/artisan websockets:serve
autostart=true
autorestart=true
user=laravel
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/websockets.log
stopwaitsecs=3600